# uwsgi/uwsgi-daemon.conf
# copy this file to /etc/init

description "harmony uWSGI"
start on runlevel [2345]
stop on runlevel [06]
respawn

exec uwsgi /var/www/harmony-solution/uwsgi/uwsgi.ini